# luci-app-cifs2-mount/root/etc/uci-defaults/luci-cifs2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@cifs2[-1]
	add ucitrack cifs2
	set ucitrack.@cifs2[-1].init=cifs2
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0